# agent/rootfs/etc/s6-overlay/s6-rc.d/init-setup/up
/bin/bash -c "set -e; mkdir -p /var/log/claworc; chmod 755 /var/log/claworc; test -f /home/claworc/.bashrc || cp -a /etc/skel/. /home/claworc/; mkdir -p /home/claworc/chrome-data /home/claworc/.openclaw; chown -R claworc:claworc /home/claworc"
